# image/rootfs/etc/systemd/system/rendermatic-update.timer
[Unit]
Description=Rendermatic OTA update timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
