# image/rootfs/etc/systemd/system/rendermatic-update.service
[Unit]
Description=Rendermatic OTA update check
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/ota-update.sh
