# image/rootfs/etc/init.d/hostname-gen
#!/sbin/openrc-run

name="hostname-gen"
description="Generate unique hostname from MAC address"

depend() {
    before networking hostname data-mount
    after modules
}

start() {
    ebegin "Generating hostname from MAC address"

    # Find first non-loopback MAC address
    local mac=""
    for iface in /sys/class/net/*; do
        [ -e "$iface/address" ] || continue
        local name=$(basename "$iface")
        [ "$name" = "lo" ] && continue
        mac=$(cat "$iface/address")
        [ -n "$mac" ] && [ "$mac" != "00:00:00:00:00:00" ] && break
        mac=""
    done

    if [ -n "$mac" ]; then
        # Last 3 bytes of MAC, e.g. 08:00:27:ca:6b:f5 → ca6bf5
        local suffix=$(echo "$mac" | awk -F: '{printf "%s%s%s", $4, $5, $6}')
        local newhost="rendermatic-${suffix}"
        hostname "$newhost"
        einfo "Hostname set to $newhost"
    else
        hostname "rendermatic"
        ewarn "No MAC address found, using default hostname"
    fi

    eend 0
}
