# image/rootfs/etc/init.d/rendermatic
#!/sbin/openrc-run

name="rendermatic"
description="Rendermatic Graphics Application"
command="/rendermatic/rendermatic"
directory="/rendermatic"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"
output_log="/data/logs/rendermatic.log"
error_log="/data/logs/rendermatic.log"

depend() {
    need localmount data-mount
    after networking dbus avahi-daemon
}

start_pre() {
    if [ ! -e /dev/fb0 ]; then
        eerror "Framebuffer device /dev/fb0 not found"
        return 1
    fi
    checkpath --directory /data/logs
}
